# sys-apps/apparmor/files/apparmor-init
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Load all configured profiles for the AppArmor security module."
description_reload="Reload all profiles"

extra_started_commands="reload"


start() {
	ebegin "Starting AppArmor"
	eindent
	apparmor_load.sh

	eoutdent
}

stop() {
	ebegin "Stopping AppArmor"
	eindent
	apparmor_unload.sh

	eoutdent
}

reload() {
	# todo: split out clean_profiles into its own function upstream
	# so we can do parse_profiles reload && clean_profiles
	# and do a proper reload instead of restart
	apparmor_restart
}
